# bin/systemd/gunicorn.service
# place in /etc/systemd/system/
# confirm the path in `WorkingDirectory` and username in `User`

[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=live
Group=www-data
WorkingDirectory=/home/live/app
ExecStart=/home/live/miniconda3/envs/edc/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --timeout=300 \
    --bind unix:/run/gunicorn.sock \
    ambition_edc.wsgi:application

[Install]
WantedBy=multi-user.target
