# bin/systemd/gunicorn-uat.socket
[Unit]
Description=gunicorn socket (UAT)

[Socket]
ListenStream=/run/gunicorn-uat.sock

[Install]
WantedBy=sockets.target